# vm_creation/raft_live_daemon.service
[Unit]
Description=raft start setup vm

[Service]
Type=forking
ExecStart=/bin/sh ls
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
