# src/raft_discovery.service
[Unit]
Description=raft bootstrap node

[Service]
Type=forking
ExecStart=/bin/sh /root/raft_project_runtime/nmap_discover.sh
RemainAfterExit=no
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target



